# supervisord-dashing.conf
[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:dashboard]
command=dashing start
directory=/dashboard
autorestart=true
restartasgroup=true
